# assets/nvidia-zerodb.service
[Unit]
Description=0dB service for NVIDIA

[Service]
Type=simple
ExecStart=/usr/bin/nvidia-zerodb

[Install]
WantedBy=multi-user.target
